# http_server.conf
[program:http]
command=/usr/bin/http-server /home/web/www -p 80
autostart=true
autorestart=true
startsecs=5