# mysql.conf
[program:mysql]
user=mysql
directory=/usr/local/
command=service mysql start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/mysql.log
stderr_logfile=/tmp/mysql.err